# upstart/redwolf_control.conf
author "RedWolf Security Inc."
description "redwolf_control is a UI and API server to manage RedWolf agents"

start on (restart_redwolf_control
          or start_redwolf_control)

stop on (restart_redwolf_control
         or stop_redwolf_control
         or stopping dbus
         or shutdown )

respawn
respawn limit 5 60 # Give up after 5 respawns in 60 seconds

pre-start script
    message="upstart: redwolf_control.conf (date `date`): pre-start script"
    logger -p local1.notice -t redwolf_control.conf.info -- ${message}
    echo $message >> /var/log/redwolf_control.log
end script

exec node /opt/control/control.js >> /var/log/redwolf_control.log 2>&1

post-stop script
    message="upstart: redwolf_control.conf (date `date`): post-stop script"
    logger -p local1.notice -t redwolf_control.conf.info -- ${message}
    echo $message >> /var/log/redwolf_control.log
end script
